# installer/overlay/etc/init.d/headless
#!/sbin/openrc-run

description="My first init script"
name="Headless init"

command="/usr/local/bin/headless"
command_background=false
pidfile="/run/${RC_SVCNAME}.pid"
