# installer/templates/wendy.service.j2
[Unit]
Description=Wendy API service unit file.

[Service]
ExecStart=/opt/bob/wendy/start-wendy.sh

[Install]
WantedBy=multi-user.target
