# installer/templates/wendy.service.j2
[Unit]
Description=Wendy Web service unit file for BoB.

[Service]
ExecStart=/opt/bob/wendy/start-wendy.sh

[Install]
WantedBy=multi-user.target
